# legacy/modules/branches/5.1/app-postfix/app-filter-training.cron
*/5 * * * *  root /usr/sbin/filtertraining >/dev/null 2>&1
