# packages/clearos-base/trunk/etc/cron.d/app-servicewatch
*/5 * * * * root LANG=en_US /usr/share/clearos/base/servicewatch >/dev/null
