# legacy/modules/branches/5.2/app-awstats/app-awstats.cron
15 6 * * * root /usr/sbin/awstatsreports
