# legacy/modules/branches/5.2/app-kolab/app-kolab.cron
#!/bin/sh

if [ -e /var/lib/imap/annotations.db ]; then
	cp /var/lib/imap/annotations.db /var/lib/imap/annotations.db.`date '+%A'`
fi
