# legacy/modules/branches/5.1/app-remote-backup/app-remote-backup.cron
0 0 * * * root /var/webconfig/scripts/rbs-schedule.php >/dev/null 2>&1
