# webconfig/apps/disk_usage/trunk/packaging/app-disk-usage.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
HOME=/usr/share/philesight

30 4 * * * root ./philesight --db /var/clearos/framework/tmp/ps.db --index / >/dev/null 2>&1
