# legacy/modules/branches/5.1/app-setup/etc/cron.d/app-servicewatch
*/5 * * * * root LANG=en_US /usr/share/system/scripts/servicewatch >/dev/null
