# legacy/modules/tags/5.2.1/app-squid/app-squid-report.cron
1 * * * *  root /usr/sbin/log2mysql >/dev/null 2>&1
