# legacy/modules/tags/5.2.1/app-squid/app-squid-purge.cron
40 3 * * *  root /usr/sbin/proxy-report-purge >/dev/null 2>&1
