# legacy/modules/tags/5.2.1/app-spamassassin/app-spamassassin.cron
43 10 * * 3 root /usr/sbin/app-sa-update >/dev/null 2>&1
