# legacy/modules/trunk/app-flexshare-api/app-flexshare.cron
*/5 * * * * root /var/flexshare/scripts/check_mbox.php >/dev/null 2>&1
