# legacy/modules/branches/5.2/app-snort/app-snort-report.cron
37 * * * * root LANG=en_US nice -n 19 /usr/share/system/scripts/snortreport >/dev/null
