# packages/clearos-console/tags/6.1.0.beta2-1/clearos-console.conf
start on stopped rc RUNLEVEL=[2345]
stop on runlevel [016]

respawn
instance tty1
exec /usr/bin/openvt -fwc 1 -- /bin/login -f clearconsole
