# legacy/modules/branches/5.2/app-raid/55-mptctl.rules
ACTION!="add", GOTO="mptctl_end"
SUBSYSTEM=="pci", DRIVER=="mpt*", ENV{MODALIAS}=="?*", RUN+="/sbin/modprobe mptctl"
LABEL="mptctl_end"
