# legacy/modules/branches/5.1/app-cyrus-imapd/app-cyrus-imapd.cron
30 5 * * *  cyrus /usr/lib/cyrus-imapd/quota -f >/dev/null 2>&1
45 5 * * *  cyrus /usr/lib/cyrus-imapd/squatter -r -v user >/dev/null 2>&1
